# ops/systemd/tanaghum-postgres-backup.service
[Unit]
Description=Tanaghum PostgreSQL Docker backup
Wants=docker.service
After=docker.service
OnFailure=tanaghum-operations-alert@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/tanaghum-hybrid
Environment=DATABASE_BACKUP_DIR=/var/backups/tanaghum/postgres
Environment=POSTGRES_CONTAINER=tanaghum-postgres-1
Environment=POSTGRES_USER=tanaghum
Environment=POSTGRES_DB=tanaghum
ExecStart=/bin/bash /opt/tanaghum-hybrid/scripts/backup-postgres-docker.sh
