# ops/systemd/tanaghum-postgres-backup.timer
[Unit]
Description=Run Tanaghum PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
